# launchd/com.barry.codexbar-serve.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  服務一：codexbar serve（資料源頭）。
  安裝時 cp 到 ~/Library/LaunchAgents/，並將 __REPLACE_WITH_DASHBOARD_TOKEN__ 換成
  openssl rand -hex 32 產生的實際 token。此檔進版控時 token 必為 placeholder。
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.barry.codexbar-serve</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/codexbar</string>
        <string>serve</string>
        <string>--host</string>
        <string>127.0.0.1</string>
        <string>--port</string>
        <string>8080</string>
        <string>--refresh-interval</string>
        <string>60</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>CODEXBAR_DASHBOARD_TOKEN</key>
        <string>__REPLACE_WITH_DASHBOARD_TOKEN__</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/Documents/vibe-coding/ai-status-dashboard/logs/codexbar-serve.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/Documents/vibe-coding/ai-status-dashboard/logs/codexbar-serve.err.log</string>
</dict>
</plist>
