# launchd/com.barry.ai-status-dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  服務二：我方 Node dashboard server（畫面）。
  node 路徑寫死絕對路徑（launchd 不載入 nvm shim）；換 nvm 版本時需更新此路徑（見 README 維護提醒）。
  安裝時 cp 到 ~/Library/LaunchAgents/，並將 __REPLACE_WITH_DASHBOARD_TOKEN__ 換成實際 token。
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.barry.ai-status-dashboard</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/YOUR_USERNAME/.nvm/versions/node/v22.17.0/bin/node</string>
        <string>/Users/YOUR_USERNAME/Documents/vibe-coding/ai-status-dashboard/src/server.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/YOUR_USERNAME/Documents/vibe-coding/ai-status-dashboard</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PORT</key>
        <string>8787</string>
        <!-- v0.3.0 起：資料源為 tokscale CLI。launchd 的 PATH 不含 nvm/homebrew，
             且 tokscale 為 node 腳本（shebang env node），故 PATH 需含 node 所在的
             nvm bin 目錄，否則子行程會噴 `env: node: No such file or directory`。
             用 `dirname "$(command -v tokscale)"` 取得該目錄後替換下方。
             TOKSCALE_BIN 為絕對路徑保險（tokscale.js 支援此覆寫）。 -->
        <key>PATH</key>
        <string>/Users/YOUR_USERNAME/.nvm/versions/node/vXX/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>TOKSCALE_BIN</key>
        <string>/Users/YOUR_USERNAME/.nvm/versions/node/vXX/bin/tokscale</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/Documents/vibe-coding/ai-status-dashboard/logs/out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/Documents/vibe-coding/ai-status-dashboard/logs/err.log</string>
</dict>
</plist>
